# crontab
37 10 * * 1,4 root python /app/main.py >> /var/log/cron.log 2>&1
