# crontab
25 20 * * 1,4 root python /app/main.py >> /var/log/cron.log 2>&1
