# crontab
15 16 * * 1,4 python /app/main.py >> /var/log/cron.log 2>&1


